# systemd/litemode-autostart.service
[Unit]
Description=LiteMode Auto-Start Services
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=opc
WorkingDirectory=/home/opc/litemode
Environment="SHELL=/bin/bash"
Environment="HOME=/home/opc"
Environment="PATH=/home/opc/.local/bin:/home/opc/bin:/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin"
Environment="NODE_OPTIONS=--max-old-space-size=2048"
Environment="DASHBOARD_PASSWORD=pantos10"
ExecStart=/bin/bash /home/opc/litemode/scripts/start-all.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
